# mysql-8.0.30-1/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/mysql.conf
etc/ld.so.conf.d/
etc/ld.so.conf.d/mysql.conf
etc/logrotate.d/
etc/logrotate.d/mysqlrouter
etc/mysql/
etc/mysql/my.cnf
etc/mysqlrouter/
etc/mysqlrouter/mysqlrouter.conf
run/
usr/
usr/bin/
usr/bin/comp_err
usr/bin/ibd2sdi
usr/bin/innochecksum
usr/bin/my_print_defaults
usr/bin/myisam_ftdump
usr/bin/myisamchk
usr/bin/myisamlog
usr/bin/myisampack
usr/bin/mysql_keyring_encryption_test
usr/bin/mysql_migrate_keyring
usr/bin/mysql_secure_installation
usr/bin/mysql_ssl_rsa_setup
usr/bin/mysql_tzinfo_to_sql
usr/bin/mysql_upgrade
usr/bin/mysqlbinlog
usr/bin/mysqld
usr/bin/mysqld_pre_systemd
usr/bin/mysqld_safe
usr/bin/mysqldumpslow
usr/bin/mysqlpump
usr/bin/mysqlrouter
usr/bin/mysqlrouter_keyring
usr/bin/mysqlrouter_passwd
usr/bin/mysqlrouter_plugin_info
usr/bin/mysqltest
usr/bin/mysqltest_safe_process
usr/bin/mysqlxtest
usr/bin/perror
usr/bin/zlib_decompress
usr/lib/
usr/lib/mysql/
usr/lib/mysql/plugin/
usr/lib/mysql/plugin/adt_null.so
usr/lib/mysql/plugin/auth.so
usr/lib/mysql/plugin/auth_socket.so
usr/lib/mysql/plugin/auth_test_plugin.so
usr/lib/mysql/plugin/component_audit_api_message_emit.so
usr/lib/mysql/plugin/component_example_component1.so
usr/lib/mysql/plugin/component_example_component2.so
usr/lib/mysql/plugin/component_example_component3.so
usr/lib/mysql/plugin/component_keyring_file.so
usr/lib/mysql/plugin/component_log_filter_dragnet.so
usr/lib/mysql/plugin/component_log_sink_json.so
usr/lib/mysql/plugin/component_log_sink_syseventlog.so
usr/lib/mysql/plugin/component_log_sink_test.so
usr/lib/mysql/plugin/component_mysqlbackup.so
usr/lib/mysql/plugin/component_mysqlx_global_reset.so
usr/lib/mysql/plugin/component_pfs_example.so
usr/lib/mysql/plugin/component_pfs_example_component_population.so
usr/lib/mysql/plugin/component_query_attributes.so
usr/lib/mysql/plugin/component_reference_cache.so
usr/lib/mysql/plugin/component_test_audit_api_message.so
usr/lib/mysql/plugin/component_test_backup_lock_service.so
usr/lib/mysql/plugin/component_test_component_deinit.so
usr/lib/mysql/plugin/component_test_host_application_signal.so
usr/lib/mysql/plugin/component_test_mysql_current_thread_reader.so
usr/lib/mysql/plugin/component_test_mysql_runtime_error.so
usr/lib/mysql/plugin/component_test_mysql_system_variable_set.so
usr/lib/mysql/plugin/component_test_pfs_notification.so
usr/lib/mysql/plugin/component_test_pfs_resource_group.so
usr/lib/mysql/plugin/component_test_sensitive_system_variables.so
usr/lib/mysql/plugin/component_test_status_var_service.so
usr/lib/mysql/plugin/component_test_status_var_service_int.so
usr/lib/mysql/plugin/component_test_status_var_service_reg_only.so
usr/lib/mysql/plugin/component_test_status_var_service_str.so
usr/lib/mysql/plugin/component_test_status_var_service_unreg_only.so
usr/lib/mysql/plugin/component_test_string_service.so
usr/lib/mysql/plugin/component_test_string_service_charset.so
usr/lib/mysql/plugin/component_test_string_service_long.so
usr/lib/mysql/plugin/component_test_sys_var_service.so
usr/lib/mysql/plugin/component_test_sys_var_service_int.so
usr/lib/mysql/plugin/component_test_sys_var_service_same.so
usr/lib/mysql/plugin/component_test_sys_var_service_str.so
usr/lib/mysql/plugin/component_test_system_variable_source.so
usr/lib/mysql/plugin/component_test_table_access.so
usr/lib/mysql/plugin/component_test_udf_registration.so
usr/lib/mysql/plugin/component_test_udf_services.so
usr/lib/mysql/plugin/component_udf_reg_3_func.so
usr/lib/mysql/plugin/component_udf_reg_avg_func.so
usr/lib/mysql/plugin/component_udf_reg_int_func.so
usr/lib/mysql/plugin/component_udf_reg_int_same_func.so
usr/lib/mysql/plugin/component_udf_reg_only_3_func.so
usr/lib/mysql/plugin/component_udf_reg_real_func.so
usr/lib/mysql/plugin/component_udf_unreg_3_func.so
usr/lib/mysql/plugin/component_udf_unreg_int_func.so
usr/lib/mysql/plugin/component_udf_unreg_real_func.so
usr/lib/mysql/plugin/component_validate_password.so
usr/lib/mysql/plugin/conflicting_variables.so
usr/lib/mysql/plugin/connection_control.so
usr/lib/mysql/plugin/daemon_example.ini
usr/lib/mysql/plugin/ddl_rewriter.so
usr/lib/mysql/plugin/debug/
usr/lib/mysql/plugin/group_replication.so
usr/lib/mysql/plugin/ha_mock.so
usr/lib/mysql/plugin/keyring_file.so
usr/lib/mysql/plugin/keyring_udf.so
usr/lib/mysql/plugin/libdaemon_example.so
usr/lib/mysql/plugin/libtest_framework.so
usr/lib/mysql/plugin/libtest_services.so
usr/lib/mysql/plugin/libtest_services_threaded.so
usr/lib/mysql/plugin/libtest_session_attach.so
usr/lib/mysql/plugin/libtest_session_detach.so
usr/lib/mysql/plugin/libtest_session_in_thd.so
usr/lib/mysql/plugin/libtest_session_info.so
usr/lib/mysql/plugin/libtest_sql_2_sessions.so
usr/lib/mysql/plugin/libtest_sql_all_col_types.so
usr/lib/mysql/plugin/libtest_sql_cmds_1.so
usr/lib/mysql/plugin/libtest_sql_commit.so
usr/lib/mysql/plugin/libtest_sql_complex.so
usr/lib/mysql/plugin/libtest_sql_errors.so
usr/lib/mysql/plugin/libtest_sql_lock.so
usr/lib/mysql/plugin/libtest_sql_processlist.so
usr/lib/mysql/plugin/libtest_sql_replication.so
usr/lib/mysql/plugin/libtest_sql_reset_connection.so
usr/lib/mysql/plugin/libtest_sql_shutdown.so
usr/lib/mysql/plugin/libtest_sql_sleep_is_connected.so
usr/lib/mysql/plugin/libtest_sql_sqlmode.so
usr/lib/mysql/plugin/libtest_sql_stmt.so
usr/lib/mysql/plugin/libtest_sql_stored_procedures_functions.so
usr/lib/mysql/plugin/libtest_sql_views_triggers.so
usr/lib/mysql/plugin/libtest_x_sessions_deinit.so
usr/lib/mysql/plugin/libtest_x_sessions_init.so
usr/lib/mysql/plugin/locking_service.so
usr/lib/mysql/plugin/mypluglib.so
usr/lib/mysql/plugin/mysql_clone.so
usr/lib/mysql/plugin/mysql_no_login.so
usr/lib/mysql/plugin/pfs_example_plugin_employee.so
usr/lib/mysql/plugin/qa_auth_client.so
usr/lib/mysql/plugin/qa_auth_interface.so
usr/lib/mysql/plugin/qa_auth_server.so
usr/lib/mysql/plugin/replication_observers_example_plugin.so
usr/lib/mysql/plugin/rewrite_example.so
usr/lib/mysql/plugin/rewriter.so
usr/lib/mysql/plugin/semisync_master.so
usr/lib/mysql/plugin/semisync_replica.so
usr/lib/mysql/plugin/semisync_slave.so
usr/lib/mysql/plugin/semisync_source.so
usr/lib/mysql/plugin/test_security_context.so
usr/lib/mysql/plugin/test_services_host_application_signal.so
usr/lib/mysql/plugin/test_services_plugin_registry.so
usr/lib/mysql/plugin/test_udf_services.so
usr/lib/mysql/plugin/udf_example.so
usr/lib/mysql/plugin/validate_password.so
usr/lib/mysql/plugin/version_token.so
usr/lib/mysqlrouter/
usr/lib/mysqlrouter/connection_pool.so
usr/lib/mysqlrouter/http_auth_backend.so
usr/lib/mysqlrouter/http_auth_realm.so
usr/lib/mysqlrouter/http_server.so
usr/lib/mysqlrouter/io.so
usr/lib/mysqlrouter/keepalive.so
usr/lib/mysqlrouter/metadata_cache.so
usr/lib/mysqlrouter/private/
usr/lib/mysqlrouter/private/libmysqlharness.so.1
usr/lib/mysqlrouter/private/libmysqlharness_stdx.so.1
usr/lib/mysqlrouter/private/libmysqlharness_tls.so.1
usr/lib/mysqlrouter/private/libmysqlrouter.so.1
usr/lib/mysqlrouter/private/libmysqlrouter_connection_pool.so.1
usr/lib/mysqlrouter/private/libmysqlrouter_http.so.1
usr/lib/mysqlrouter/private/libmysqlrouter_http_auth_backend.so.1
usr/lib/mysqlrouter/private/libmysqlrouter_http_auth_realm.so.1
usr/lib/mysqlrouter/private/libmysqlrouter_io_component.so.1
usr/lib/mysqlrouter/private/libmysqlrouter_metadata_cache.so.1
usr/lib/mysqlrouter/private/libmysqlrouter_mysqlxmessages.so.1
usr/lib/mysqlrouter/private/libmysqlrouter_routing.so.1
usr/lib/mysqlrouter/private/libprotobuf-lite.so.3.19.4
usr/lib/mysqlrouter/rest_api.so
usr/lib/mysqlrouter/rest_connection_pool.so
usr/lib/mysqlrouter/rest_metadata_cache.so
usr/lib/mysqlrouter/rest_router.so
usr/lib/mysqlrouter/rest_routing.so
usr/lib/mysqlrouter/router_openssl.so
usr/lib/mysqlrouter/router_protobuf.so
usr/lib/mysqlrouter/routing.so
usr/lib/private/
usr/lib/private/libprotobuf-lite.so.3.19.4
usr/lib/private/libprotobuf.so.3.19.4
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/mysqld.service
usr/lib/systemd/system/mysqld@.service
usr/lib/systemd/system/mysqlrouter.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/mysql.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/mysql.conf
usr/lib/tmpfiles.d/mysqlrouter.conf
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/comp_err.1.gz
usr/share/man/man1/ibd2sdi.1.gz
usr/share/man/man1/innochecksum.1.gz
usr/share/man/man1/lz4_decompress.1.gz
usr/share/man/man1/my_print_defaults.1.gz
usr/share/man/man1/myisam_ftdump.1.gz
usr/share/man/man1/myisamchk.1.gz
usr/share/man/man1/myisamlog.1.gz
usr/share/man/man1/myisampack.1.gz
usr/share/man/man1/mysql_secure_installation.1.gz
usr/share/man/man1/mysql_ssl_rsa_setup.1.gz
usr/share/man/man1/mysql_tzinfo_to_sql.1.gz
usr/share/man/man1/mysql_upgrade.1.gz
usr/share/man/man1/mysqlbinlog.1.gz
usr/share/man/man1/mysqldumpslow.1.gz
usr/share/man/man1/mysqlman.1.gz
usr/share/man/man1/mysqlpump.1.gz
usr/share/man/man1/mysqlrouter.1.gz
usr/share/man/man1/mysqlrouter_passwd.1.gz
usr/share/man/man1/mysqlrouter_plugin_info.1.gz
usr/share/man/man1/perror.1.gz
usr/share/man/man1/zlib_decompress.1.gz
usr/share/man/man8/
usr/share/man/man8/mysqld.8.gz
usr/share/mysql/
usr/share/mysql/LICENSE
usr/share/mysql/LICENSE-test
usr/share/mysql/README
usr/share/mysql/README-test
usr/share/mysql/bulgarian/
usr/share/mysql/bulgarian/errmsg.sys
usr/share/mysql/charsets/
usr/share/mysql/charsets/Index.xml
usr/share/mysql/charsets/README
usr/share/mysql/charsets/armscii8.xml
usr/share/mysql/charsets/ascii.xml
usr/share/mysql/charsets/cp1250.xml
usr/share/mysql/charsets/cp1251.xml
usr/share/mysql/charsets/cp1256.xml
usr/share/mysql/charsets/cp1257.xml
usr/share/mysql/charsets/cp850.xml
usr/share/mysql/charsets/cp852.xml
usr/share/mysql/charsets/cp866.xml
usr/share/mysql/charsets/dec8.xml
usr/share/mysql/charsets/geostd8.xml
usr/share/mysql/charsets/greek.xml
usr/share/mysql/charsets/hebrew.xml
usr/share/mysql/charsets/hp8.xml
usr/share/mysql/charsets/keybcs2.xml
usr/share/mysql/charsets/koi8r.xml
usr/share/mysql/charsets/koi8u.xml
usr/share/mysql/charsets/latin1.xml
usr/share/mysql/charsets/latin2.xml
usr/share/mysql/charsets/latin5.xml
usr/share/mysql/charsets/latin7.xml
usr/share/mysql/charsets/macce.xml
usr/share/mysql/charsets/macroman.xml
usr/share/mysql/charsets/swe7.xml
usr/share/mysql/czech/
usr/share/mysql/czech/errmsg.sys
usr/share/mysql/danish/
usr/share/mysql/danish/errmsg.sys
usr/share/mysql/dictionary.txt
usr/share/mysql/docs/
usr/share/mysql/docs/ChangeLog
usr/share/mysql/docs/INFO_BIN
usr/share/mysql/docs/INFO_SRC
usr/share/mysql/docs/LICENSE.router
usr/share/mysql/docs/README.router
usr/share/mysql/docs/sample_mysqlrouter.conf
usr/share/mysql/dutch/
usr/share/mysql/dutch/errmsg.sys
usr/share/mysql/english/
usr/share/mysql/english/errmsg.sys
usr/share/mysql/estonian/
usr/share/mysql/estonian/errmsg.sys
usr/share/mysql/french/
usr/share/mysql/french/errmsg.sys
usr/share/mysql/german/
usr/share/mysql/german/errmsg.sys
usr/share/mysql/greek/
usr/share/mysql/greek/errmsg.sys
usr/share/mysql/hungarian/
usr/share/mysql/hungarian/errmsg.sys
usr/share/mysql/install_rewriter.sql
usr/share/mysql/italian/
usr/share/mysql/italian/errmsg.sys
usr/share/mysql/japanese/
usr/share/mysql/japanese/errmsg.sys
usr/share/mysql/korean/
usr/share/mysql/korean/errmsg.sys
usr/share/mysql/messages_to_clients.txt
usr/share/mysql/messages_to_error_log.txt
usr/share/mysql/mysql-log-rotate
usr/share/mysql/norwegian-ny/
usr/share/mysql/norwegian-ny/errmsg.sys
usr/share/mysql/norwegian/
usr/share/mysql/norwegian/errmsg.sys
usr/share/mysql/polish/
usr/share/mysql/polish/errmsg.sys
usr/share/mysql/portuguese/
usr/share/mysql/portuguese/errmsg.sys
usr/share/mysql/romanian/
usr/share/mysql/romanian/errmsg.sys
usr/share/mysql/russian/
usr/share/mysql/russian/errmsg.sys
usr/share/mysql/serbian/
usr/share/mysql/serbian/errmsg.sys
usr/share/mysql/slovak/
usr/share/mysql/slovak/errmsg.sys
usr/share/mysql/spanish/
usr/share/mysql/spanish/errmsg.sys
usr/share/mysql/swedish/
usr/share/mysql/swedish/errmsg.sys
usr/share/mysql/ukrainian/
usr/share/mysql/ukrainian/errmsg.sys
usr/share/mysql/uninstall_rewriter.sql
var/
var/lib/
var/lib/mysqlrouter/
var/log/
var/log/mysqlrouter/
